# systemd/plantwaterer.timer
[Unit]
Description=Run plant watering every two weeks

[Timer]
OnUnitActiveSec=14d
Persistent=true

[Install]
WantedBy=timers.target
